# etc/supervisord.d/es-register.ini
[program:es-register]
command=/opt/qnib/elasticsearch/bin/register.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_events_enabled = false
stderr_events_enabled = false
autorestart = false
autostart = true
startsecs=0
startretries=50
exitcodes=0
